# deploy/systemd/shopping-mcp.service
[Unit]
Description=Drission Shopping MCP
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/drission-shopping-mcp
EnvironmentFile=/home/pi/drission-shopping-mcp/.env
ExecStart=/home/pi/drission-shopping-mcp/.venv/bin/python -m shopping_mcp.asgi
Restart=always
RestartSec=5

# --- Hardening ---
# Limit blast radius if the Chromium process or the app itself is exploited.
# Chromium needs /tmp for its own sandbox bits, so PrivateTmp=true is fine.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
# Allow writing only to the app dir (for debug_captures) and the chromium
# user-data-dir. Adjust paths if your layout differs from DEPLOY.md.
ReadWritePaths=/home/pi/drission-shopping-mcp /home/pi/.cache/drission-shopping-mcp

[Install]
WantedBy=multi-user.target
